# build/copyfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
#command=/usr/sbin/sshd
autostart=true
autorestart=true

[program:slapd]
command=/usr/sbin/slapd
autostart=true
autorestart=true

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND
autostart=true
autorestart=true
