# build/copyfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
#command=/usr/sbin/sshd -D
command=service ssh start
#command=/usr/sbin/sshd
#autostart=true
a#utorestart=true

[program:slapd]
#command=/usr/sbin/slapd
#command=/usr/sbin/slapd -h ldap:/// ldapi:/// -g openldap -u openldap -F /etc/ldap/slapd.d
command=service slapd start
#autostart=true
#autorestart=true

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND
autostart=true
autorestart=true
